# layers/Persistence.json
{
  "name": "Elastic Security: Persistence",
  "versions": {
    "attack": "8",
    "navigator": "4.2",
    "layer": "4.1"
  },
  "domain": "enterprise-attack",
  "description": "Mapping Elastic Secrity (https://www.elastic.co/security) Detection Rules to MITRE ATT&CK Navigator",
  "filters": {},
  "sorting": 0,
  "layout": {
    "layout": "side",
    "showName": true,
    "showID": true
  },
  "hideDisabled": false,
  "legendItems": [],
  "gradient": {
    "colors": [
      "#bbffdd",
      "#44ff66"
    ],
    "minValue": 0,
    "maxValue": 47
  },
  "showTacticRowBackground": true,
  "tacticRowBackground": "#bbddff",
  "selectTechniquesAcrossTactics": false,
  "selectSubtechniquesWithParent": false,
  "metadata": [
    {
      "name": "Elastic version",
      "value": "7.12.0"
    }
  ],
  "techniques": [
    {
      "techniqueID": "T1546",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / high",
          "value": "Potential Modification of Accessibility Binaries (73)"
        },
        {
          "name": "EQL / high",
          "value": "Potential Modification of Accessibility Binaries (73)"
        },
        {
          "name": "Query / low",
          "value": "Potential Application Shimming via Sdbinst (21)"
        },
        {
          "name": "Query / low",
          "value": "Potential Application Shimming via Sdbinst (21)"
        },
        {
          "name": "EQL / medium",
          "value": "Installation of Custom Shim Databases (21)"
        },
        {
          "name": "EQL / medium",
          "value": "Image File Execution Options Injection (41)"
        },
        {
          "name": "EQL / medium",
          "value": "Registry Persistence via AppCert DLL (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Registry Persistence via AppInit DLL (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Component Object Model Hijacking (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Component Object Model Hijacking (47)"
        },
        {
          "name": "EQL / low",
          "value": "Persistence via WMI Event Subscription (21)"
        }
      ],
      "score": 11
    },
    {
      "techniqueID": "T1136",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / low",
          "value": "User Account Creation (21)"
        },
        {
          "name": "EQL / high",
          "value": "Creation of a Hidden Local User Account (73)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1574",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / low",
          "value": "Adobe Hijack Persistence (21)"
        },
        {
          "name": "Query / low",
          "value": "Adobe Hijack Persistence (21)"
        },
        {
          "name": "EQL / high",
          "value": "Suspicious DLL Loaded for Persistence or Privilege Escalation (73)"
        },
        {
          "name": "EQL / high",
          "value": "Suspicious DLL Loaded for Persistence or Privilege Escalation (73)"
        }
      ],
      "score": 4
    },
    {
      "techniqueID": "T1505",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / medium",
          "value": "Potential Shell via Web Server (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1053",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / low",
          "value": "Local Scheduled Task Commands (21)"
        },
        {
          "name": "Query / high",
          "value": "Persistence via TelemetryController Scheduled Task Hijack (73)"
        },
        {
          "name": "Query / low",
          "value": "Creation or Modification of a new GPO Scheduled Task or Service (21)"
        },
        {
          "name": "EQL / low",
          "value": "Suspicious Image Load (taskschd.dll) from MS Office (21)"
        },
        {
          "name": "EQL / medium",
          "value": "Suspicious Execution via Scheduled Task (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Suspicious Execution via Scheduled Task (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Scheduled Task Created by a Windows Script (47)"
        }
      ],
      "score": 7
    },
    {
      "techniqueID": "T1543",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / medium",
          "value": "System Shells via Services (47)"
        },
        {
          "name": "EQL / high",
          "value": "Persistence via Update Orchestrator Service Hijack (73)"
        },
        {
          "name": "EQL / high",
          "value": "Suspicious ImagePath Service Creation (73)"
        },
        {
          "name": "EQL / low",
          "value": "Unusual Persistence via Services Registry (21)"
        },
        {
          "name": "EQL / high",
          "value": "Suspicious ImagePath Service Creation (73)"
        },
        {
          "name": "EQL / low",
          "value": "Unusual Persistence via Services Registry (21)"
        },
        {
          "name": "EQL / low",
          "value": "LaunchDaemon Creation or Modification and Immediate Loading (21)"
        },
        {
          "name": "EQL / low",
          "value": "Launch Agent Creation or Modification and Immediate Loading (21)"
        },
        {
          "name": "Query / medium",
          "value": "Persistence via Docker Shortcut Modification (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Finder Sync Plugin Registered and Enabled (47)"
        },
        {
          "name": "Query / medium",
          "value": "Modification of Standard Authentication Module or Configuration (47)"
        },
        {
          "name": "Query / medium",
          "value": "Modification of OpenSSH Binaries (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Creation of Hidden Launch Agent or Daemon (47)"
        },
        {
          "name": "Query / medium",
          "value": "Suspicious Hidden Child Process of Launchd (47)"
        }
      ],
      "score": 14
    },
    {
      "techniqueID": "T1547",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / low",
          "value": "Persistence via Kernel Module Modification (21)"
        },
        {
          "name": "EQL / high",
          "value": "Persistence via Hidden Run Key Detected (73)"
        },
        {
          "name": "EQL / medium",
          "value": "Installation of Security Support Provider (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Shortcut File Written or Modified for Persistence (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Persistent Scripts in the Startup Directory (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Installation of Security Support Provider (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Shortcut File Written or Modified for Persistence (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Persistent Scripts in the Startup Directory (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Execution of Persistent Suspicious Program (47)"
        },
        {
          "name": "EQL / low",
          "value": "Startup or Run Key Registry Modification (21)"
        },
        {
          "name": "EQL / medium",
          "value": "Startup Folder Persistence via Unsigned Process (41)"
        },
        {
          "name": "Query / medium",
          "value": "Potential Persistence via Login Hook (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Creation of Hidden Login Item via Apple Script (47)"
        }
      ],
      "score": 13
    },
    {
      "techniqueID": "T1137",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / high",
          "value": "Persistence via Microsoft Office AddIns (73)"
        },
        {
          "name": "EQL / medium",
          "value": "Persistence via Microsoft Outlook VBA (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Persistence via Microsoft Outlook VBA (47)"
        }
      ],
      "score": 3
    },
    {
      "techniqueID": "T1112",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / medium",
          "value": "Uncommon Registry Persistence Change (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1037",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / medium",
          "value": "Persistence via Folder Action Script (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Persistence via Login or Logout Hook (47)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1059",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / medium",
          "value": "Persistence via Folder Action Script (47)"
        },
        {
          "name": "EQL / medium",
          "value": "Creation of Hidden Login Item via Apple Script (47)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1556",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / high",
          "value": "Potential OpenSSH Backdoor Logging Activity (73)"
        },
        {
          "name": "Query / medium",
          "value": "Modification of Standard Authentication Module or Configuration (47)"
        },
        {
          "name": "Query / medium",
          "value": "Modification of OpenSSH Binaries (47)"
        }
      ],
      "score": 3
    },
    {
      "techniqueID": "T1554",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / high",
          "value": "Potential OpenSSH Backdoor Logging Activity (73)"
        },
        {
          "name": "EQL / low",
          "value": "Sublime Plugin or Application Script Modification (21)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1098",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / medium",
          "value": "SSH Authorized Keys File Modification (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1078",
      "enabled": true,
      "metadata": [
        {
          "name": "Query / medium",
          "value": "Potential Hidden Local User Account Creation (47)"
        },
        {
          "name": "Query / medium",
          "value": "Attempt to Enable the Root Account (47)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1564",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / medium",
          "value": "Creation of Hidden Launch Agent or Daemon (47)"
        },
        {
          "name": "Query / medium",
          "value": "Suspicious Hidden Child Process of Launchd (47)"
        }
      ],
      "score": 2
    }
  ]
}