# layers/7.13.1/Lateral Movement.json
{
  "name": "El.Sec. Lateral Movement",
  "versions": {
    "attack": "9",
    "navigator": "4.3",
    "layer": "4.2"
  },
  "domain": "enterprise-attack",
  "description": "Mapping Elastic Secrity (https://www.elastic.co/security) Detection Rules to MITRE ATT&CK Navigator",
  "filters": {},
  "sorting": 0,
  "layout": {
    "layout": "side",
    "showName": true,
    "showID": true
  },
  "hideDisabled": false,
  "legendItems": [
    {
      "label": "El.Sec. => Elastic Secrity",
      "color": "#bbddff"
    },
    {
      "label": "EQL => Event Query Language",
      "color": "#fa744e"
    },
    {
      "label": "ML => Machine Learning",
      "color": "#0077cc"
    },
    {
      "label": "Q => Query (KQL, Lucene, Elasticsearch DSL)",
      "color": "#54bcb2"
    },
    {
      "label": "TH => Threshold",
      "color": "#fec514"
    }
  ],
  "gradient": {
    "colors": [
      "#bbffdd",
      "#44ff66"
    ],
    "minValue": 0,
    "maxValue": 42
  },
  "showTacticRowBackground": true,
  "tacticRowBackground": "#bbddff",
  "selectTechniquesAcrossTactics": false,
  "selectSubtechniquesWithParent": false,
  "metadata": [
    {
      "name": "Elastic Stack version",
      "value": "7.13.1"
    },
    {
      "name": "Metadata format",
      "value": "<Type> / <Severity>: <Name> (<Risk_score>)"
    },
    {
      "name": "Github project",
      "value": "https://github.com/ElasticSA/elsec_dr2an"
    }
  ],
  "techniques": [
    {
      "techniqueID": "T1021",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Virtual Private Network Connection Attempt (21)"
        },
        {
          "name": "EQL / Med.",
          "value": "NullSessionPipe Registry Modification (47)"
        },
        {
          "name": "EQL / High",
          "value": "Potential Remote Desktop Shadowing Activity (73)"
        },
        {
          "name": "EQL / Low",
          "value": "Virtual Private Network Connection Attempt (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Med.",
          "value": "Connection to External Network via Telnet (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Mounting Hidden or WebDav Remote Shares (21)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remote Execution via File Shares (47)"
        },
        {
          "name": "EQL / High",
          "value": "Incoming DCOM Lateral Movement via MSHTA (73)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remote File Copy to a Hidden Share (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Connection to Internal Network via Telnet (47)"
        },
        {
          "name": "EQL / High",
          "value": "Execution via TSClient Mountpoint (73)"
        },
        {
          "name": "EQL / Med.",
          "value": "Incoming DCOM Lateral Movement with ShellBrowserWindow or ShellWindows (47)"
        },
        {
          "name": "EQL / High",
          "value": "Incoming DCOM Lateral Movement with MMC (73)"
        },
        {
          "name": "EQL / High",
          "value": "Potential SharpRDP Behavior (73)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remotely Started Services via RPC (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remote Scheduled Task Creation (47)"
        },
        {
          "name": "Q / Med.",
          "value": "Remote SSH Login Enabled via systemsetup Command (47)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Control Spawned via Script Interpreter (21)"
        },
        {
          "name": "EQL / High",
          "value": "Potential Remote Desktop Tunneling Detected (73)"
        },
        {
          "name": "EQL / High",
          "value": "Lateral Movement via Startup Folder (73)"
        }
      ],
      "score": 22
    },
    {
      "techniqueID": "T1543",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1543.003",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1569",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1569.002",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        },
        {
          "name": "EQL / Low",
          "value": "Service Command Lateral Movement (21)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1550",
      "enabled": true,
      "metadata": [
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        },
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1550.003",
      "enabled": true,
      "metadata": [
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        },
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1558",
      "enabled": true,
      "metadata": [
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        },
        {
          "name": "Q / High",
          "value": "Potential Kerberos Attack via Bifrost (73)"
        }
      ],
      "score": 2
    },
    {
      "techniqueID": "T1210",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Med.",
          "value": "Direct Outbound SMB Connection (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Direct Outbound SMB Connection (47)"
        },
        {
          "name": "Q / Med.",
          "value": "Abnormally Large DNS Response (47)"
        }
      ],
      "score": 3
    },
    {
      "techniqueID": "T1021.002",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Med.",
          "value": "NullSessionPipe Registry Modification (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Mounting Hidden or WebDav Remote Shares (21)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remote Execution via File Shares (47)"
        },
        {
          "name": "EQL / Med.",
          "value": "Remote File Copy to a Hidden Share (47)"
        }
      ],
      "score": 4
    },
    {
      "techniqueID": "T1570",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Med.",
          "value": "Lateral Tool Transfer (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1047",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Med.",
          "value": "WMI Incoming Lateral Movement (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1053",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / Med.",
          "value": "Remote Scheduled Task Creation (47)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1547",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / High",
          "value": "Lateral Movement via Startup Folder (73)"
        }
      ],
      "score": 1
    },
    {
      "techniqueID": "T1547.001",
      "enabled": true,
      "metadata": [
        {
          "name": "EQL / High",
          "value": "Lateral Movement via Startup Folder (73)"
        }
      ],
      "score": 1
    }
  ]
}